# vm/service/tasklist.service
[Unit]
Description=Tasklist Spring Boot Application
After=network.target
Wants=postgresql.service

[Service]
Type=simple
User=siseko
Group=siseko
ExecStart=/usr/bin/java -jar /opt/tasklist/app/tasklist-api.jar
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10

# Environment variables for Spring Boot placeholders
Environment=DB_URL=jdbc:postgresql://192.168.18.3:5432/tasklistdb
Environment=DB_USERNAME=postgres
Environment=DB_PASSWORD=admin
Environment=JPA_DDL_AUTO=update
Environment=JPA_SHOW_SQL=true
Environment=SERVER_PORT=8080
Environment=LOG_LEVEL_SPRING=INFO
Environment=LOG_LEVEL_TASKLIST=DEBUG
Environment=LOG_FILE=/opt/tasklist/logs/tasklist.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tasklist/logs /tmp

# Resource limits
LimitNOFILE=65536
MemoryLimit=512M

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
